# etc/cron.d/thachanpy-notification
SHELL=/bin/sh
0 6 * * * root /data/thachanpy/thachanpy-notification/thachanpy-notification
